# supervisor-trp-api.conf
[unix_http_server]
file=/var/run/supervisor.sock

[supervisord]
logfile=/var/log/supervisor/supervisord.log
pidfile=/var/run/supervisord.pid
childlogdir=/var/log/supervisor

[rpcinterface:supervisor]
supervisor.rpcinterface_factory = supervisor.rpcinterface:make_main_rpcinterface

[supervisorctl]
serverurl=unix:///var/run/supervisor.sock

[program:trp-api]
directory=/opt/trp-api
command=/opt/trp-api/venv/bin/gunicorn \
    --config /opt/trp-api/gunicorn_config.py \
    app_basic_auth:app
user=ubuntu
autostart=true
autorestart=true
startsecs=10
stopwaitsecs=10
stdout_logfile=/var/log/trp-api/gunicorn.log
stdout_logfile_maxbytes=100MB
stdout_logfile_backups=10
stderr_logfile=/var/log/trp-api/gunicorn_error.log
stderr_logfile_maxbytes=100MB
stderr_logfile_backups=10
environment=PATH="/opt/trp-api/venv/bin",HOME="/home/ubuntu"
